# etc/docker/node/supervisord.conf.d/00-webpack.conf
[program:webpack-dev-server]
command=/usr/local/bin/npm run dev-server
directory=/var/www/app
autostart=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[group:webpack]
programs=webpack-dev-server
