# etc/docker/php/supervisord.conf
[supervisord]
logfile=/var/log/supervisor/supervisord.log
loglevel=warn
pidfile=/run/supervisor/supervisord.pid
nodaemon=true

[unix_http_server]
file=/run/supervisor/supervisord.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/supervisor/supervisord.sock

[program:php-fpm]
command=/usr/local/sbin/php-fpm
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[include]
files=/etc/supervisord.conf.d/*.conf
